# reboot-checker.timer
[Unit]
Description=Notify on server reboot

[Timer]
OnBootSec=60s
AccuracySec=60s

[Install]
WantedBy=timers.target

