# torban.timer
[Unit]
Description=Ban IPs from the tor network

[Timer]
OnBootSec=2min
OnUnitActiveSec=6h
AccuracySec=5min

[Install]
WantedBy=timers.target